# scripts/launchd/event-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ai.dollars-and-cents.event-refresh
  Pulls Federal Register events daily before the factory tick so that
  PREMISE_INVALIDATION alerts have fresh inputs.

  Schedule: 06:00 ET = 03:00 PT, Mon-Fri (mini's system TZ is PT;
  launchd fires in system-local time. See CLAUDE.md DST notes).

  Deploy: sed -i '' "s|\$HOME|$HOME|g" then move to ~/Library/LaunchAgents/.
  See OPERATIONS.md.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.dollars-and-cents.event-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>$HOME/.cents/bin/cents-wrap</string>
    <string>event-refresh</string>
    <string>event</string>
    <string>refresh</string>
  </array>

  <!-- 06:00 ET = 03:00 PT weekdays. launchd fires in system-local TZ;
       wrapper's TZ env affects process time but NOT firing time. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>$HOME/.cents/logs/launchd-event-refresh.out</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.cents/logs/launchd-event-refresh.err</string>
</dict>
</plist>
