# scripts/launchd/shadow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ai.dollars-and-cents.shadow
  Backfill 30-day forward returns for symbols the factory considered but
  rejected. Builds the shadow_opens table — counterfactual labeled
  outcomes for "what if we had opened these?" analytics. Doesn't affect
  the live experiment cohort.

  Schedule: 07:00 ET, Mon-Fri (after both factory arms have ticked).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.dollars-and-cents.shadow</string>

  <key>ProgramArguments</key>
  <array>
    <string>$HOME/.cents/bin/cents-wrap</string>
    <string>shadow</string>
    <string>shadow</string>
    <string>backfill</string>
    <string>--horizon</string>
    <string>30</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>$HOME/.cents/logs/launchd-shadow.out</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.cents/logs/launchd-shadow.err</string>
</dict>
</plist>
